# v5/.claude/PAI/PULSE/MenuBar/com.pai.pulse-menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.pulse-menubar</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/Applications/PAI Pulse.app/Contents/MacOS/PAI Pulse</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PAI_PULSE_DIR</key>
        <string>__HOME__/.claude/PAI/PULSE</string>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/.claude/PAI/PULSE/logs/menubar-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/PAI/PULSE/logs/menubar-stderr.log</string>
</dict>
</plist>
